# daemon-init-redhat.in
#!/bin/sh
#
# Icinga	Startup script for the Icinga monitoring daemon
#
# chkconfig:	- 85 15
# description:	Icinga is a service monitoring system
### BEGIN INIT INFO
# Provides:		icinga
# Required-Start:	$local_fs $syslog $network
# Default-Start:	3 4 5
# Default-Stop: 	0 1 6
# Short-Description:    start and stop Icinga monitoring server
# Description:		Icinga is is a service monitoring system 
### END INIT INFO

# Source function library.
. /etc/rc.d/init.d/functions

prefix="@prefix@"
exec_prefix="@exec_prefix@"
exec="@bindir@/icinga"
prog="icinga"
config="@sysconfdir@/icinga.cfg"
pidfile="@lockfile@"
user="@icinga_user@"

[ -e /etc/sysconfig/$prog ] && . /etc/sysconfig/$prog

lockfile=/var/lock/subsys/$prog

start() {
    [ -x $exec ] || exit 5
    [ -f $config ] || exit 6
    echo -n $"Starting $prog: "
    daemon --user=$user $exec -d $config
    retval=$?
    echo
    [ $retval -eq 0 ] && touch $lockfile
    return $retval
}

stop() {
    echo -n $"Stopping $prog: "
    killproc -p ${pidfile} -d 10 $exec
    retval=$?
    echo
    [ $retval -eq 0 ] && rm -f $lockfile
    return $retval
}


restart() {
    stop
    start
}

reload() {
    echo -n $"Reloading $prog: "
    killproc -p ${pidfile} $exec -HUP
    RETVAL=$?
    echo
}

force_reload() {
    restart
}

check_config() {
        $nice runuser -s /bin/bash - $user -c "$corelimit >/dev/null 2>&1 ; $exec -v $config > /dev/null 2>&1"
	RETVAL=$?
        if [ $RETVAL -ne 0 ] ; then
                echo -n $"Configuration validation failed"
                failure
                echo
                exit 1
        fi
}

case "$1" in
    start)
        status $prog && exit 0
	check_config
        $1
        ;;
    stop)
        status $prog|| exit 0
        $1
        ;;
    restart)
	check_config
        $1
        ;;
    reload)
        status $prog || exit 7
	check_config
        $1
        ;;
    force-reload)
        force_reload
        ;;
    status)
        status $prog
        ;;
    condrestart|try-restart)
        status $prog|| exit 0
        restart
        ;;
   configtest)
        echo -n  $"Checking config for $prog: "
        check_config && success
        echo

	;;
    *)
        echo $"Usage: $0 {start|stop|status|restart|condrestart|try-restart|reload|force-reload|configtest}"
        exit 2
esac
exit $?

